# crontab
0 3 * * * $HOME/kadira-compose/flush-mongo
